# asam-deb-build/opt/asam-linux-installer/asam.service
[Unit]
Description=ASAM Launcher Service
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple

# Run ASAM as a normal user (change 'jez' to whichever user should run ASAM)
User=jez
Group=jez

# Working directory for ASAM
WorkingDirectory=/opt/asam

# Launcher script
ExecStart=/opt/asam-linux-installer/asam_systemd_launcher.sh

# Restart policy
Restart=on-failure
RestartSec=5

# Environment for Wine
Environment=WINEDEBUG=-all
Environment=WINEPREFIX=/opt/asam-wine-prefix

# Logging
StandardOutput=append:/var/log/asam_launcher.log
StandardError=append:/var/log/asam_launcher.log

# Ensure Wine has access to a display if needed
Environment=DISPLAY=:0

[Install]
WantedBy=multi-user.target